# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/base-project/bin/uwsgi --ini /srv/base-project/.config/uwsgi/uwsgi-app.ini

